# debian/postinst
#!/bin/bash
#
# @license  http://www.gnu.org/licenses/gpl.html GPL Version 3
# @author   OpenMediaVault Plugin Developers <plugins@omv-extras.org>
# @copyright Copyright (c) 2013-2019 OpenMediaVault Plugin Developers
#
# This program is free software: you can redistribute it and/or modify
# it under the terms of the GNU General Public License as published by
# the Free Software Foundation, either version 3 of the License, or
# any later version.
#
# This program is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
# GNU General Public License for more details.
#
# You should have received a copy of the GNU General Public License
# along with this program. If not, see <http://www.gnu.org/licenses/>.

set -e

. /etc/default/openmediavault
. /usr/share/openmediavault/scripts/helper-functions

case "${1}" in
    configure)

        if [ -z "${2}" ]; then
            deb-systemd-helper disable autoshutdown.service >/dev/null || :
            deb-systemd-invoke stop autoshutdown.service >/dev/null || :
        fi

        # Activate package triggers.
        dpkg-trigger update-fixperms || :
        dpkg-trigger update-locale || :

        # Initialize and migrate configuration database.
        echo "Updating configuration database ..."
        omv-confdbadm create "conf.service.autoshutdown"
        if [ -n "${2}" ]; then
            echo "Migrating configuration database ..."
            omv-confdbadm migrate "conf.service.autoshutdown" "${2}"
        fi

        if dpkg --compare-versions "${2}" lt-nl "5.1.7"; then
            echo "Removing /etc/pm/sleep.d/autoshutdown-restart"
            rm -f "/etc/pm/sleep.d/autoshutdown-restart"
        fi

        if dpkg --compare-versions "${2}" lt-nl "5.1.8"; then
            echo "Removing /usr/lib/systemd/system-sleep/autoshutdown-restart"
            rm -f "/usr/lib/systemd/system-sleep/autoshutdown-restart"
        fi

        # Deploy the configuration for autoshutdown services.
        echo "Deploying autoshutdown service configuration"
        LANG=C.UTF-8 omv-salt deploy run --no-color --quiet autoshutdown
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument '${1}'" >&2
        exit 1
    ;;
esac

#DEBHELPER#

exit 0
